# mt7902-late.service
# /etc/systemd/system/mt7902-late.service
#
# MT7902 WiFi Driver Early Initialization
# Loads before network services to ensure WiFi is available at boot

[Unit]
Description=MT7902 WiFi Driver Initialization
Documentation=https://github.com/your-repo/gen4-mt7902
# Start BEFORE network services try to come up
Before=network-pre.target
Wants=network-pre.target
# Wait for basic system initialization
After=systemd-udev-settle.service
After=systemd-modules-load.service
# Don't pull in standard dependencies
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
# Run as root (required for insmod)
User=root
# Set working directory to driver location
WorkingDirectory=/home/dan/builds/gen4-mt7902
# Allow time for driver initialization
TimeoutStartSec=60
# Small delay for PCIe stability
ExecStartPre=/usr/bin/sleep 3
# Load the driver
ExecStart=/usr/bin/bash /home/dan/builds/gen4-mt7902/load.sh
# Verify the interface appeared
ExecStartPost=/bin/bash -c '\
  echo "Verifying WiFi interface appeared..."; \
  for i in {1..20}; do \
    if ip link show 2>/dev/null | grep -q wlan; then \
      WLAN_IF=$(ip link show | grep -oP "wlan\\d+" | head -1); \
      echo "✓ WiFi interface ready: $WLAN_IF"; \
      exit 0; \
    fi; \
    sleep 0.5; \
  done; \
  echo "ERROR: WiFi interface did not appear within 10 seconds" >&2; \
  exit 1'
# Give NetworkManager a moment to detect the new interface
ExecStartPost=/bin/bash -c '\
  echo "Waiting for NetworkManager..."; \
  sleep 2; \
  if command -v nmcli >/dev/null 2>&1; then \
    nmcli general reload 2>/dev/null || true; \
  fi; \
  echo "✓ Driver initialization complete"'
# On stop, unload the module cleanly
ExecStop=/bin/bash -c '\
  if lsmod | grep -q "^mt7902\\b"; then \
    echo "Unloading mt7902 module..."; \
    rmmod mt7902 || rmmod -f mt7902 || true; \
  fi'
# Restart policy
Restart=no
# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt7902-wifi

[Install]
# Start during early boot
WantedBy=network-pre.target
WantedBy=multi-user.target