# mt7902-late.service
# /etc/systemd/system/mt7902-late.service
[Unit]
Description=MT7902 WiFi Driver — Late load
Documentation=https://github.com/your-repo/gen4-mt7902
After=graphical.target systemd-user-sessions.service
Before=iwd.service
Wants=iwd.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
WorkingDirectory=/home/dan/builds/gen4-mt7902
TimeoutStartSec=180

ExecStartPre=/usr/bin/sleep 15
ExecStartPre=/bin/echo "[mt7902] late-load starting" | /usr/bin/systemd-cat -t mt7902 -p info

ExecStart=/bin/bash /home/dan/builds/gen4-mt7902/load.sh
ExecStop=/bin/bash -c '/sbin/rmmod mt7902 2>/dev/null || true; echo "[mt7902] rmmod attempted" | /usr/bin/systemd-cat -t mt7902 -p info'

Restart=no

[Install]
WantedBy=graphical.target
